# Packaging/deb/postinst
#!/bin/sh
set -e

BASE="/var/mobile/Library/Application Support/EasyTier"
PLIST="/var/jb/Library/LaunchDaemons/com.zeroninx.easytierd.plist"
APP="/var/jb/Applications/EasyTier.app"
PKG_LOG="$BASE/logs/package.log"

mkdir -p "$BASE/profiles" "$BASE/logs" "$BASE/runtime" "$BASE/exports"
chown -R mobile:mobile "$BASE" 2>/dev/null || true
chmod 700 "$BASE" "$BASE/profiles" "$BASE/logs" "$BASE/runtime" "$BASE/exports"

log() {
    echo "$(date '+%Y-%m-%d %H:%M:%S') postinst: $*" >> "$PKG_LOG"
}

run_logged() {
    log "$*"
    "$@" >> "$PKG_LOG" 2>&1 || log "failed: $*"
}

log "starting"

if [ -x /usr/bin/uicache ]; then
    run_logged /usr/bin/uicache -p "$APP"
elif [ -x /var/jb/usr/bin/uicache ]; then
    run_logged /var/jb/usr/bin/uicache -p "$APP"
fi

if [ -x /bin/launchctl ] && [ -f "$PLIST" ]; then
    run_logged /bin/launchctl bootout system "$PLIST"
    if [ -x /usr/bin/killall ]; then
        run_logged /usr/bin/killall easytierd
    elif [ -x /var/jb/usr/bin/killall ]; then
        run_logged /var/jb/usr/bin/killall easytierd
    fi
    run_logged /bin/launchctl bootstrap system "$PLIST"
    run_logged /bin/launchctl enable system/com.zeroninx.easytierd
    run_logged /bin/launchctl kickstart -k system/com.zeroninx.easytierd
    run_logged /bin/launchctl print system/com.zeroninx.easytierd
    run_logged /bin/launchctl load -w "$PLIST"
fi

log "finished"

exit 0
